# crontab
1 * * * * python3 /PycharmProjects/UCSF_Assignment_AP3/main.py
