# crontab
*/1 * * * * python3 main.py >> /tmp/cron.out
